# deploy/com.pantoken.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchDaemon TEMPLATE for the pantoken server. A system LaunchDaemon (not a per-user
  LaunchAgent) so it starts at boot with nobody logged in — the Mac Mini runs and
  reboots headlessly. It runs as your user via <UserName> so pi keeps your perms/env.

  Don't install this file directly — the @@…@@ placeholders are rendered by
  `scripts/deploy-ctl.sh install` (which sudo-copies the result to
  /Library/LaunchDaemons). Secrets are NOT here: run.sh sources them from
  ~/.local/state/pantoken/pantoken.env. See deploy/DEPLOY.md.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pantoken.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@LIVE@@/deploy/run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@@LIVE@@</string>

  <!-- Run as you, not root: lets the deploy poller restart us with a plain kill. -->
  <key>UserName</key>
  <string>@@USER@@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- A boot-time daemon has no login shell: set HOME + PATH explicitly so bun is
         found and ~/.ssh, ~/.local/state etc. resolve. -->
    <key>HOME</key>
    <string>@@HOME@@</string>
    <key>PATH</key>
    <string>@@HOME@@/.bun/bin:@@HOME@@/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>@@LOGDIR@@/pantoken.out.log</string>
  <key>StandardErrorPath</key>
  <string>@@LOGDIR@@/pantoken.err.log</string>
</dict>
</plist>
